# web/entrypoint.conf
[program:entrypoint]
command=/usr/bin/entrypoint.sh
autostart=true
autorestart=unexpected
exitcodes=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr