# web/letsencrypt.conf
[program:letsencrypt]
command=/usr/bin/letsencrypt.sh
autostart=true
autorestart=false