# web/nginx.conf
[program:nginx]
command=/usr/bin/nginx.sh
autostart=true
autorestart=true